# profile.d/_brew.bash
# shellcheck shell=bash
export PATH=/opt/homebrew/opt/gnu-sed/libexec/gnubin:/opt/homebrew/opt/gnu-getopt/bin:/opt/homebrew/bin:/opt/homebrew/opt/node@14/bin:/Users/pkrohn/bin:$PATH

# asdf needs `brew` in the path
. $(brew --prefix asdf)/libexec/asdf.sh
