# profile.d/aws.bash
export AWS_PROFILE=dev1-cdp002
export AWS_REGION=us-west-2

if [ -f /opt/homebrew/etc/bash_completion ]; then
  . /opt/homebrew/etc/bash_completion
fi
