# profile.d/aws.bash
export AWS_PROFILE=dev1-cdp002
export AWS_REGION=us-west-2
