# profile.d/brew.bash
export PATH=/opt/homebrew/opt/gnu-sed/libexec/gnubin:/opt/homebrew/opt/gnu-getopt/bin:/opt/homebrew/bin:/opt/homebrew/opt/node@14/bin:/Users/pkrohn/bin:$PATH
