# profile.d/cdp.bash
# for the vmf cli wrapper
export SPINNAKER_REPO="${HOME}/source/spinnaker"
