# profile.d/git.bash

if [ -f /opt/homebrew/etc/bash_completion.d/git-completion.bash ]; then
    # shellcheck source=/opt/homebrew/etc/bash_completion.d/git-completion.bash
  . /opt/homebrew/etc/bash_completion.d/git-completion.bash
fi
