# profile.d/java.bash
export JAVA_HOME=/usr/local/openjdk_11.0.19_11.64.20_aarch64
export PATH=${JAVA_HOME}/bin:${PATH}
